# scripts/cron.daily/hrstka-backup.sh
#!/bin/bash
"$HRSTKA_DEPLOYMENT_DIR/scripts/backup.sh" $HRSTKA_MONGO_HOST $HRSTKA_MONGO_PORT $HRSTKA_MONGO_DB $HRSTKA_BACKUP_DIR